# debian/postinst
#! /bin/sh
# postinst script for w3c-markup-validator
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the `postinst' is called with `abort-upgrade',
#     `abort-remove' or `abort-deconfigure'.

case "$1" in
    configure|reconfigure)

        server="apache2"

        # Old versions stored config files here
        # We won't delete the old config file in case there are
        # customizations. However we must delete the link pointing to it
        # so that the upgrade can proceed.
        for old in apache2 apache apache-ssl apache-perl; do
            if [ -s /etc/$old/w3c-markup-validator.conf ]; then
                rm /etc/$old/w3c-markup-validator.conf
            fi
        done

	if [ ! -s /etc/$server/conf.d/w3c-markup-validator.conf ]; then
	    ln -s /etc/w3c/httpd.conf   \
                    /etc/$server/conf.d/w3c-markup-validator.conf
        fi

        if [ $server = "apache2" ]; then

            if [ ! -s /etc/$server/mods-enabled/include.load ]; then
                a2enmod include
            fi

            if [ ! -s /etc/$server/mods-enabled/rewrite.load ]; then
                a2enmod rewrite
            fi
	    
	fi

	if [ -x "/etc/init.d/$server" ]; then
		if [ -x /usr/sbin/invoke-rc.d ]; then
			invoke-rc.d $server restart || true
		else
			/etc/init.d/$server restart || true
		fi
	fi
	sleep 2

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


